# systemd/wai-celery-beat.service
[Unit]
Description=WAI Telegram Celery Beat (Scheduler)
After=network.target docker.service wai-celery.service
Requires=docker.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/wai-telegram/packages/backend
Environment="PATH=/root/.local/bin:/usr/local/bin:/usr/bin"
EnvironmentFile=/opt/wai-telegram/.env.production

ExecStart=/root/.local/bin/uv run celery -A app.tasks.celery_app:celery_app beat --loglevel=info
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID

Restart=always
RestartSec=10
StartLimitBurst=3
StartLimitIntervalSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wai-celery-beat

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/wai-telegram /tmp

[Install]
WantedBy=multi-user.target
